# vm-monitor.service
[Unit]
Description=Azure VM Scheduled Events Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/monitor_vm_scheduled_events.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 